# etc/profile.d/tomcat7.sh
NAME="tomcat7"

export TOMCAT_HOME="/var/lib/$NAME"
export CATALINA_HOME="/usr/share/$NAME"
export CATALINA_BASE="/var/lib/$NAME"
